# root/etc/s6-overlay/s6-rc.d/sshd/run
#!/command/with-contenv bash

exec 2>&1 \
    s6-notifyoncheck -d -n 300 -w 1000 -c "nc -z localhost 22" \
        /usr/sbin/sshd -D -e -p 22
